# assets/service/shadowsocks/run
#!/bin/sh
ss-local -c /root/assets/config.json
